# flight-tracker-user.service
[Unit]
Description=Flight Tracker Server (User Service)
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/pi/berrybase_demos/flightaware_demo
Environment="HOME=/home/pi"
Environment="USER=pi"
Environment="PATH=/home/pi/flightaware-venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONUNBUFFERED=1"
ExecStart=/bin/bash /home/pi/berrybase_demos/flightaware_demo/start_flight_tracker.sh
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target

